# keboola-config/supervisord/services/node.conf
[program:nextjs-app]
command=/usr/local/bin/node .next/standalone/server.js
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=10
environment=PATH="/usr/local/bin:/usr/bin:/bin",NODE_ENV="production",PORT="3000",HOSTNAME="0.0.0.0"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
